# fiddupay.service
[Unit]
Description=FidduPay Cryptocurrency Payment Gateway
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=vibes
WorkingDirectory=/home/vibes/crypto-payment-gateway
Environment="RUST_LOG=info"
EnvironmentFile=/home/vibes/crypto-payment-gateway/.env.production
ExecStart=/home/vibes/crypto-payment-gateway/target/release/crypto-payment-gateway
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
